# containers/minetest-server/root/etc/s6-overlay/s6-rc.d/init-voxelibre/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

echo "Setting up VoxeLibre (MineClone2) game..."

# Check if VoxeLibre is already installed
if [ -d "/config/.minetest/games/VoxeLibre" ]; then
    echo "VoxeLibre already installed, skipping download"
    lsiown -R abc:abc /config/.minetest/games/VoxeLibre
    exit 0
fi

echo "Downloading VoxeLibre from GitHub..."

# Create temporary directory
mkdir -p /tmp/voxelibre

# Download latest VoxeLibre release
VOXELIBRE_VERSION=$(curl -sX GET "https://api.github.com/repos/VoxeLibre/VoxeLibre/releases/latest" | jq -r .tag_name)
echo "Downloading VoxeLibre version: $VOXELIBRE_VERSION"

curl -o /tmp/voxelibre.tar.gz -L \
    "https://github.com/VoxeLibre/VoxeLibre/archive/${VOXELIBRE_VERSION}.tar.gz"

# Extract to games directory
tar xf /tmp/voxelibre.tar.gz -C /tmp/voxelibre --strip-components=1
mv /tmp/voxelibre /config/.minetest/games/VoxeLibre

# Set proper ownership
lsiown -R abc:abc /config/.minetest/games/VoxeLibre

# Update minetest.conf to use VoxeLibre as default game
if [ -f "/config/.minetest/main-config/minetest.conf" ]; then
    # Remove existing default_game line if present
    sed -i '/^default_game/d' /config/.minetest/main-config/minetest.conf
    
    # Add VoxeLibre as default game
    echo "default_game = VoxeLibre" >> /config/.minetest/main-config/minetest.conf
    
    # Add admin user "Treffy" with all privileges
    if ! grep -q "name = Treffy" /config/.minetest/main-config/minetest.conf; then
        echo "" >> /config/.minetest/main-config/minetest.conf
        echo "# Admin user configuration" >> /config/.minetest/main-config/minetest.conf
        echo "name = Treffy" >> /config/.minetest/main-config/minetest.conf
    fi
    
    echo "Minetest configuration updated for VoxeLibre"
fi

# Cleanup
rm -rf /tmp/voxelibre.tar.gz

echo "VoxeLibre installation complete"