# containers/victim/wazuh-install.service
[Unit]
Description=Install Wazuh Agent on First Boot
After=multi-user.target
ConditionPathExists=!/var/ossec/.installed

[Service]
Type=oneshot
ExecStart=/opt/purple-team/scripts/install-wazuh.sh
EnvironmentFile=/etc/environment.wazuh
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target